# puppet/modules/systemd/templates/weather.service.erb
[Unit]
Description=Wheather collecting. 

[Service]
Type=forking
ExecStart=/usr/local/bin/weather-init.sh start

[Install]
WantedBy=multi-user.target
